# files/twitcoin-upstart.conf
description "twitcoin"

start on runlevel [2345]
stop on runlevel [016]

console log

exec start-stop-daemon --start -c twitcoin --exec /home/twitcoin/start.sh
